# custom-scripts/replicate/resources/rsync.service
[Unit]
Description="rsync-wrapper service unit file"

[Service]
ExecStart="/usr/local/bin/rsync-wrapper"
