# custom-scripts/sync-data/resources/rsync.timer
[Unit]
Description="Run rsync (wrapper) script hourly"

[Timer]
OnCalendar=*-*-* *:00:00
Unit=rsync.service

[Install]
WantedBy=multi-user.target
